# supervisor.d/wkhtml2x.ini
[program:wkhtml2x]
directory=/var/www/html
command=/usr/local/bin/php wkhtml2x.php
user=root
autorestart=true
startretries=5
stderr_logfile=/var/log/wkhtml2x/err.log
stdout_logfile=/var/log/wkhtml2x/out.log