# supervisor.d/wkhtml2x.ini
[program:wkhtml2x]
command=/usr/local/bin/php wkhtml2x.php
directory=/var/www/html
user=root
stopsignal=INT
autostart=true
autorestart=false
stderr_logfile=/var/log/wkhtml2x/err.log
stdout_logfile=/var/log/wkhtml2x/out.log